# debian/unrealircd.postinst
#!/bin/sh
# postinst script for unrealircd
#
# see: dh_installdeb(1)

chgrp irc /usr/bin/ircd
chmod 750 /usr/bin/ircd

touch /etc/unrealircd/ircd.log
chown -R irc:irc /etc/unrealircd
update-rc.d unrealircd defaults

exit 0
